# build/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command=php-fpm7 --nodaemonize --fpm-config /etc/php7/conf.d/php-fpm.conf -c /etc/php7/php.ini
user=www-data
priority=1
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
#exitcodes=0,2
stopsignal=QUIT
stopasgroup=false
killasgroup=false
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/php5-fpm.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_events_enabled=false
serverurl=AUTO